# pi_project/systemd/koten2026.service
[Unit]
Description=Koten2026 Pi Hand Sender
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/koten2026
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/koten2026/.venv/bin/python /home/pi/koten2026/app/pi_hand_sender.py --config /home/pi/koten2026/config/endpoint.json --model /home/pi/koten2026/hand_landmarker.task
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
